# bi_scripts/superhero2/luigi/superhero2.cron
0 2 * * * hadoop cd /home/data/bi_scripts && luigi --module superhero2\.luigi.daily_trigger RunAllTasks
